# packaging/voisu.service
[Unit]
Description=Voisu dictation daemon
After=dbus.socket pipewire.service xdg-desktop-portal.service
Wants=dbus.socket pipewire.service xdg-desktop-portal.service
PartOf=graphical-session.target
StartLimitIntervalSec=30s
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/bin/voisu-daemon --systemd
Restart=on-failure
RestartSec=2s
TimeoutStopSec=60s

# Sandboxing (defense-in-depth for a compromised dependency). Every exception
# below exists for a reason — see docs/packaging-fedora.md before tightening.
NoNewPrivileges=yes
# Read-only filesystem except the daemon's real write surfaces. The runtime dir
# (control socket, pw-record scratch, diagnostics/history) is the only
# ReadWritePaths entry; config (config.toml, dictionary, credentials) and state
# (the RemoteDesktop restore token) are provisioned by systemd through
# ConfigurationDirectory/StateDirectory. systemd creates those
# BEFORE it sets up the mount sandbox and mounts them read-write, so a fresh home
# that never had ~/.config/voisu can no longer fail namespace setup
# (status=226/NAMESPACE) the way an explicit ReadWritePaths=%h/.config/voisu did.
# For a --user unit they resolve to $XDG_CONFIG_HOME/voisu (~/.config/voisu) and
# $XDG_STATE_HOME/voisu (~/.local/state/voisu), matching config_dir()/state_root().
ProtectSystem=strict
ConfigurationDirectory=voisu
StateDirectory=voisu
ReadWritePaths=%t
PrivateTmp=yes
# AF_UNIX: control socket, D-Bus (Secret Service, portals), PipeWire.
# AF_INET/AF_INET6: provider HTTPS/WSS via curl.
# AF_NETLINK: glibc getaddrinfo enumerates interfaces over netlink.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
# The daemon and its subprocesses (pw-record, curl, wl-copy, secret-tool) do
# no JIT; validate against a real install before assuming this holds.
MemoryDenyWriteExecute=yes

[Install]
WantedBy=graphical-session.target
